# RussianTranslation/deploy/pwg-ru-max-orchestrator.service
[Unit]
Description=PWG RU four-account Max dispatcher
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=pwg-factory
WorkingDirectory=/opt/factory/SanskritLexicography/RussianTranslation
ExecStart=/usr/bin/python3 src/pilot/max_account_orchestrator.py --db /var/lib/pwg-ru/max.sqlite recover
ExecStart=/usr/bin/python3 src/pilot/max_account_orchestrator.py --db /var/lib/pwg-ru/max.sqlite run-once
ExecStart=/usr/bin/python3 src/pilot/max_account_orchestrator.py --db /var/lib/pwg-ru/max.sqlite record-done --coordinator src/pilot/coordinator.py --cwd /opt/factory/SanskritLexicography/RussianTranslation
StandardOutput=append:/var/log/pwg-ru/orchestrator.log
StandardError=append:/var/log/pwg-ru/orchestrator.log
